# usecases/bench-sync.service
[Unit]
Description=Run bench sync

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bench sync-loop

[Install]
WantedBy=multi-user.target

